# packer-ansible-docker-systemd-master/service/motd.service
[Unit]
Description=motd service

[Service]
Type=oneshot
ExecStart=/usr/bin/python /opt/motd.py

[Install]
WantedBy=multi-user.target
